# apps/cpoe_cli/packaging/linux/systemd/cpop.socket
[Unit]
Description=CPOP Socket Activation
Documentation=man:cpop(1)
PartOf=cpop.service

[Socket]
# Unix socket for local communication
ListenStream=/run/cpop/cpop.sock
SocketMode=0660
SocketUser=cpop
SocketGroup=cpop

# Accept connections
Accept=no

# Buffer sizes
ReceiveBuffer=1M
SendBuffer=1M

# Backlog
Backlog=128

# Permissions
DirectoryMode=0750

[Install]
WantedBy=sockets.target
